# diskspinner@.timer
[Unit]
Description=Scheduler for diskspinner@%i service

[Timer]
OnBootSec=0
OnUnitActiveSec=10s

[Install]
WantedBy=multi-user.target
